# packaging/debian/maranoded.service
[Unit]
Description=Maranode AI Inference Daemon
Documentation=https://github.com/maranode/maranode
After=network.target
ConditionPathExists=/usr/bin/maranoded

[Service]
Type=simple
User=maranode
Group=maranode
ExecStart=/usr/bin/maranoded
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

# Hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
ReadWritePaths=/var/lib/maranode /var/log/maranode

[Install]
WantedBy=multi-user.target
